# template/supervisord.d/redpen.ini
[program:redpen]
command = sh -c "java -jar -Dredpen.home=${REDPEN_HOME} redpen-server.war"
autorestart = true
stdout_logfile = /var/log/supervisor/redpen-server.log
redirect_stderr = true
environment = REDPEN_HOME="<REDPEN_HOME>"
directory = /home/redpen/redpen/redpen-server/target
